# root/etc/services.d/klippy/run
#!/usr/bin/with-contenv sh

exec /opt/klipper/klippy-env/bin/python /opt/klipper/klippy/klippy.py $KLIPPER_PRINTER_FILE -l $KLIPPER_LOG_FILE
